# examples/systemd/eas-backup.service
[Unit]
Description=EAS Station Automated Backup
Documentation=https://github.com/KR8MER/eas-station
After=network.target docker.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/eas-station
ExecStart=/usr/bin/python3 /opt/eas-station/tools/backup_scheduler.py \
    --output-dir /var/backups/eas-station \
    --label scheduled \
    --log-file /var/log/eas-backup.log \
    --notify admin@example.com
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eas-backup

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups/eas-station /var/log /opt/eas-station/backups

[Install]
WantedBy=multi-user.target
