# examples/systemd/eas-backup.timer
[Unit]
Description=EAS Station Backup Timer (Daily at 2 AM)
Documentation=https://github.com/KR8MER/eas-station

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
